# ops/systemd/controller-linear-factory.timer
[Unit]
Description=Run the guarded Linear coding factory every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=20s
Persistent=true
Unit=controller-linear-factory.service

[Install]
WantedBy=timers.target
