# docs/agentic-harness-gui.service.template
[Unit]
Description=Agentic Harness GUI service
After=network.target

[Service]
Type=simple
User=<USER>
WorkingDirectory=<WORKDIR>
EnvironmentFile=-<TOKEN_ENV_FILE>
ExecStart=<EXECUTABLE> --project-dir <WORKDIR> --host 127.0.0.1 --port <PORT> --no-open
Restart=on-failure

[Install]
WantedBy=multi-user.target
